# nginx/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=service nginx start
